# chatbot/supervisor/app.conf
[program:main]
command=/home/marco_mamei/venv/bin/python3 main.py
directory=/home/marco_mamei
stderr_logfile=/home/marco_mamei/log.err
stdout_logfile=/home/marco_mamei/log.out
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true